# polkit-gnome.service
[Unit]
Description=Authentication Agent for PolicyKit
Requires=dbus.socket

[Service]
ExecStart=@LIBEXECDIR@/polkit-gnome-authentication-agent-1

